# systemd-kmonad.service
[Unit]
Description=KMonad for Attack Shark keyboard (Auto-started by udev)
# might be useful if nothing else depends on this explicitly
StopWhenUnneeded=true 

[Service]
Type=simple
ExecStart=/usr/bin/kmonad /etc/kmonad/fullkeyboard-cable.kbd
Restart=on-failure
RestartSec=5s
# KillMode=process # Ensures kmonad process itself is killed, not just the service
StandardOutput=journal
StandardError=journal

# Note: No [Install] section is needed because udev will manage it.
